# frodo_vagrant/salt/django-example.conf
description "Gunicorn running a django process."

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid vagrant
setgid vagrant
chdir /vagrant/django-example/

exec /vagrant/example-venv/bin/gunicorn --bind 0.0.0.0:8080 django_example.wsgi:application
